# .setup/linux/arch/install/post-install
#!/usr/bin/env bash

# Enabling network-related services
systemctl enable dhcpcd
systemctl enable NetworkManager
