# .assets/services/ollama-runit/run
#!/bin/sh

exec /home/ismail/_store/Vault/binaries/ollama serve 1>/dev/null
