# .assets/services/ollama-runit/run
#!/bin/sh

exec /usr/local/bin/ollama serve 1>/dev/null
