# .assets/services/ollama-runit/run
#!/bin/sh

exec /home/ismail/_store/Applications/ollama serve 1>/dev/null
