# .assets/services/echonet-runit/run
#!/bin/sh

exec /home/ismail/_store/Vault/media/echonet 1>/dev/null
