# .assets/services/runsvdir-ismail/run
#!/bin/sh

export USER=ismail
export HOME="/home/$USER"

groups="$(id -Gn "$USER" | tr ' ' ':')"
svdir="$HOME/.services"

exec chpst -u "$USER:$groups" runsvdir "$svdir"